# home/sbog/.config/autostart/xbindkeys.desktop
[Desktop Entry]
Exec=xbindkeys
Icon=system-run
Path=
Terminal=false
Type=Application
